# ops/systemd/trading-agent-daily-research.timer
[Unit]
Description=JvR Trading Agent — hourly COL tick (v3.15.2)
Documentation=https://github.com/roudjy/trading-agent docs/adr/ADR-011-v3.10-architecture.md

[Timer]
# Hourly on the top of every UTC hour. The campaign launcher is a
# single-tick CLI; if it lands in worker_busy / cooldown / no-candidates
# state it exits 0 quickly. Heavy work (the actual research subprocess)
# only runs when the policy engine selects a campaign.
OnCalendar=*-*-* *:00:00 UTC
# If the system was down at the scheduled time, fire once it comes back.
Persistent=true
Unit=trading-agent-daily-research.service

[Install]
WantedBy=timers.target
