# ops/systemd/trading-agent-daily-research.service
[Unit]
Description=JvR Trading Agent — daily default research run (v3.10)
Documentation=https://github.com/roudjy/trading-agent docs/adr/ADR-011-v3.10-architecture.md
# Only fire once the dashboard container is up. The unit fails loudly and
# the journal captures it when the container is not running.
After=docker.service network-online.target

[Service]
Type=oneshot
# Run the default daily preset inside the existing jvr_dashboard container,
# so CLI + UI + scheduler all share the SAME execution path (ADR-011 §5).
ExecStart=/usr/bin/docker exec jvr_dashboard python /app/researchctl.py run trend_equities_4h_baseline
# Never auto-run the crypto_diagnostic_1h preset from this timer; it is
# excluded_from_daily_scheduler=True. That preset must stay manual.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
